# .local/share/systemd/user/pactl-auto-switch-sink.service
[Service]
Type=oneshot
RemainAfterExit=true
StandardOutput=journal
ExecStart=pactl load-module module-switch-on-connect
[Install]
WantedBy=default.target
