# systemd/user/avfs.service
[Unit]
Description=Avfs archive auto-mount
RequiresMountsFor=%t

[Service]
Type=forking
RuntimeDirectory=%P
Environment=AVFSBASE=%t/%P
ExecStart=/usr/bin/mountavfs
ExecStop=/usr/bin/umountavfs
ProtectSystem=strict

[Install]
WantedBy=default.target
